# packer/http/tmp-hostpath_pv.automount
[Unit]
Before=local-fs.target
RequiresMountsFor=/data /tmp

[Automount]
Where=/tmp/hostpath_pv
